# daemon/sup.conf
[program:suman-daemon]
command=/Users/alexamil/WebstormProjects/oresoftware/sumanjs/suman/cli/suman-daemon.sh