# daemon/supervisor.conf
[program:suman-daemon]
command=%(suman_daemon_exec_path)